# zoneminder/supervisord/supervisord-zoneminder.conf
[program:zoneminder]
command=zmpkg.pl start
user=root
autostart=true
autorestart=false
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
