# samba/supervisord/supervisord-sssd.conf
[program:sssd]
command=sssd -i -d 4
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
